# resources/systemd/devture-nagadmin.service
[Unit]
Description=Devture Nagadmin
Requires=docker.service
After=docker.service
DefaultDependencies=no

[Service]
Type=simple
WorkingDirectory=/srv/http/nagadmin
# Adjust the `just` path to wherever it is installed (e.g. /usr/local/bin/just
# for a manually-downloaded binary). `prod` runs the production environment; a
# bare `just run`/`stop` would default to dev.
ExecStartPre=-/usr/bin/just stop prod

ExecStart=/usr/bin/just run prod

ExecStop=-/usr/bin/just stop prod

Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
