# resources/systemd/devture-nagadmin.service
[Unit]
Description=Devture Nagadmin
Requires=docker.service
After=docker.service
DefaultDependencies=no

[Service]
Type=simple
WorkingDirectory=/srv/http/nagadmin
# `prod` runs the production environment; a bare `just run`/`stop` defaults to dev.
ExecStartPre=-/usr/local/bin/just stop prod

ExecStart=/usr/local/bin/just run prod

ExecStop=-/usr/local/bin/just stop prod

Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
